# systemd/photo-frame-sync.service
[Unit]
Description=Photo Frame Sync

[Service]
Type=oneshot
ExecStart=/usr/bin/rclone sync gdrive:PhotoFrame /mnt/photo-frame/photos --fast-list
