# systemd/photo-frame-sync.timer
[Unit]
Description=Sync photos every 10 minutes

[Timer]
OnBootSec=3min
OnUnitActiveSec=24h

[Install]
WantedBy=timers.target
