# systemd/photo-frame.service
[Unit]
Description=Photo Frame (Pygame)
After=network-online.target graphical.target
Wants=network-online.target

[Service]
Type=simple
User=pi
EnvironmentFile=/home/admin/photo-frame-data/env
WorkingDirectory=/home/admin/photo-frame
# Ensure we have a DISPLAY for X11 (typical Pi desktop)
Environment=DISPLAY=:0
ExecStart=/home/admin/photo-frame/scripts/run_frame.sh
Restart=on-failure
RestartSec=2

[Install]
WantedBy=graphical.target
